# iso-seed/etc/xdg/autostart/lockedin-dashboard.desktop
[Desktop Entry]
Name=LockedIn Dashboard
Comment=LockedinOS unified productivity dashboard
Exec=/usr/bin/lockedin-dashboard --no-sandbox
Icon=preferences-desktop
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Phase=Application
X-GNOME-Autostart-Delay=2
NoDisplay=true
